# supervisord-ssh.conf
[program:ssh]
command=/start-ssh.sh
numprocs=1
autostart=true
autorestart=true
